# deploy/systemd/theeye-risk-metrics-refresh.service
[Unit]
Description=TheEyeBeta risk metrics refresh (POST risk-service compute-portfolio-metrics)
Wants=network-online.target theeye-risk-service.service
After=network-online.target theeye-risk-service.service postgresql.service

[Service]
Type=oneshot
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-risk-metrics-refresh
ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.venv/bin/python scripts/risk_metrics_refresh.py
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
