# deploy/systemd/theeye-intraday-ingest.service
[Unit]
Description=TheEyeBeta canonical 15m intraday ingestion
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env
ExecStart=/home/the-eye-beta/.local/bin/uv run python -m workers.intraday_ingestion_worker --once --run-type scheduled
TimeoutStartSec=900
RuntimeMaxSec=1200

[Install]
WantedBy=multi-user.target
