# deploy/systemd/theeye-agent-runtime.service
[Unit]
Description=TheEyeBeta agent-runtime (executor agent API)
Wants=network-online.target nats.service postgresql.service redis-server.service theeyebeta-litellm.service
After=network-online.target nats.service postgresql.service redis-server.service theeyebeta-litellm.service
OnFailure=theeye-notify-failure@%n.service

[Service]
Type=simple
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/services/agent_runtime
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-agent-runtime
Environment=PYTHONPATH=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/services/agent_runtime/src
Environment=AGENT_RUNTIME_HOST=127.0.0.1
Environment=AGENT_RUNTIME_PORT=8004
Environment=LITELLM_PROXY_URL=http://127.0.0.1:4000
ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeye-agent-runtime
NoNewPrivileges=true
PrivateTmp=true
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
