# deploy/systemd/theeye-news-bridge.timer
[Unit]
Description=Run news bridge every 30 minutes (RSS articles → market_news)

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
