# deploy/systemd/theeye-news-ingest.timer
[Unit]
Description=Poll production news feeds frequently so theeyebeta.news_articles stays fresh
Requires=theeye-news-ingest.service

[Timer]
# Every 15 minutes; URL-hash deduplication makes reruns cheap and keeps
# research/news surfaces current instead of bursty.
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
