# deploy/systemd/theeye-zinc-eod.service
[Unit]
Description=TheEyeBeta ZINC INVESTMENTS end-of-day email report
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-zinc-eod
ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.venv/bin/python scripts/zinc_fund_tracker.py --eod
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeye-zinc-eod
