# deploy/systemd/theeye-risk-service.service
[Unit]
Description=TheEyeBeta risk-service (pre-trade risk + risk_metrics writer)
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=simple
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-risk-service
ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.venv/bin/python services/risk_service/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeye-risk-service
NoNewPrivileges=true
PrivateTmp=true
TimeoutStopSec=30
KillSignal=SIGTERM
OnFailure=theeye-notify-failure@%n.service

[Install]
WantedBy=multi-user.target
