# deploy/systemd/theeye-massive-ingest.service
[Unit]
Description=TheEyeBeta canonical daily price ingestion (Massive)
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-massive-ingest
ExecStart=/home/the-eye-beta/.local/bin/uv run python -m workers.massive_ingestion_worker --run-type scheduled
TimeoutStartSec=7200
RuntimeMaxSec=9000

[Install]
WantedBy=multi-user.target
