# deploy/systemd/theeye-snapshot-packager.service
[Unit]
Description=TheEyeBeta snapshot-packager (NATS consumer + agent-ready JSON, MinIO writer)
Wants=network-online.target nats.service postgresql.service
After=network-online.target nats.service postgresql.service

[Service]
Type=simple
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/services/snapshot_packager
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env.theeye-snapshot-packager
Environment=PYTHONPATH=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/services/snapshot_packager/src
ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeye-snapshot-packager
NoNewPrivileges=true
PrivateTmp=true
TimeoutStopSec=30
OnFailure=theeye-notify-failure@%n.service

[Install]
WantedBy=multi-user.target
