# deploy/systemd/theeye-sector.timer
[Unit]
Description=Run sector aggregation after the nightly ingest + indicator sync

[Timer]
OnCalendar=Mon..Fri 22:05 UTC
Persistent=true
RandomizedDelaySec=60
Unit=theeye-sector.service

[Install]
WantedBy=timers.target
