# deploy/systemd/archived/theeyebeta-watcher.service
[Unit]
Description=TheEyeBeta Repo Watcher (auto-update + restart on push)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=the-eye-beta
Group=the-eye-beta

WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaLocal

EnvironmentFile=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaLocal/.env
Environment="BRANCH=main"
Environment="CHECK_INTERVAL_SECONDS=60"
Environment="API_PORT=8000"

ExecStart=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaLocal/scripts/watch_repo_and_restart.sh

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeyebeta-watcher

NoNewPrivileges=true

TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
