# deploy/systemd/staged/theeye-risk-service.service
# STAGED — installed in deploy/systemd/staged/ on purpose so install_systemd_units.sh
# (top-level glob only) does NOT copy it to /etc yet. Do NOT `systemctl enable --now`.
#
# risk_service is the intended writer for theeyebeta.risk_metrics (issue #4). It is
# REACTIVE: it appends a risk_metrics row only when called via its HTTP bridge
# (/v1/compute-portfolio-metrics, /v1/validate-order) or gRPC (ValidateOrder,
# ComputePortfolioMetrics). On its own it writes nothing.
#
# risk_metrics is empty because the platform holds no book, NOT because this unit was
# missing. Before activating, ALL of the following must hold (see
# docs/ops/risk-metrics-activation.md for the full checklist):
#   1. C++ risk kernels built:  make build-cpp
#      Provides zinc_native._zinc_risk; without it `import risk_service.app` fails.
#   2. >=1 row in theeyebeta.portfolios (+ positions). risk_metrics.portfolio_id is a
#      NOT NULL FK to portfolios — with an empty book nothing can be computed or inserted.
#   3. A driver calling it on a schedule: the master_orchestrator RiskMetricsScheduler
#      (set RISK_SERVICE_URL + RISK_METRICS_PORTFOLIO_IDS), or a timer hitting the bridge.
#
# To activate once prerequisites are met: move this file to deploy/systemd/ (top level),
# run sudo deploy/install_systemd_units.sh, then `systemctl enable --now theeye-risk-service`.

[Unit]
Description=TheEyeBeta risk-service (pre-trade risk + theeyebeta.risk_metrics writer)
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=simple
User=the-eye-beta
Group=the-eye-beta
WorkingDirectory=/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd
EnvironmentFile=-/home/the-eye-beta/TheEyeBeta2025/TheEyeBetaProd/.env
ExecStart=/home/the-eye-beta/.local/bin/uv run python services/risk_service/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theeye-risk-service
NoNewPrivileges=true
PrivateTmp=true
TimeoutStopSec=30
KillSignal=SIGTERM
OnFailure=theeye-notify-failure@%n.service

[Install]
WantedBy=multi-user.target
