# deploy/systemd/nats.service.d/jetstream.conf
[Service]
# Enable JetStream for audit-service and trading event durability.
ExecStart=
ExecStart=/usr/local/bin/nats-server -p 4222 -a 127.0.0.1 -js -sd /var/lib/nats/jetstream
